# iol_dashboard.service
[Unit]
Description=IOL Dashboard and Master
After=network.target

[Service]
Type=simple
ExecStart=/home/user/start_iol_dashboard.sh
Restart=always
RestartSec=5
User=user
Group=user
Environment=DISPLAY=:0
Environment=XDG_RUNTIME_DIR=/run/user/1000

[Install]
WantedBy=graphical.target

