# rotorsync.service
[Unit]
Description=Rotorsync BLE GATT Server (Bumble)
After=network.target
Wants=bluetooth.target

[Service]
Type=simple
ExecStartPre=/bin/bash -c 'HCI=$(hciconfig -a | grep -B1 "00:01:95:9C:C7:F5" | head -1 | cut -d: -f1); [ -n "$HCI" ] && hciconfig $HCI down || true'
ExecStart=/usr/bin/python3 /opt/rotorsync_bumble.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
