# rotorsync_watchdog.service
[Unit]
Description=Rotorsync BLE Watchdog
After=rotorsync.service
Requires=rotorsync.service

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/rotorsync_watchdog.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
