# systemd/rotorlink.service
[Unit]
Description=RotorLink local WiFi link (iPad <-> dashboard bridge)
# Bridges the dashboard command socket — but connects at runtime with retry, so it
# must NOT order After=iol_dashboard.service: that service is pulled into graphical.target
# (ordered after multi-user.target), while this unit is WantedBy=multi-user.target, which
# creates an ordering cycle that systemd breaks by DROPPING iol_dashboard at boot (dashboard
# never starts). Order only on the network.
After=network-online.target NetworkManager.service
Wants=network-online.target

[Service]
Type=simple
# root: config writes over WiFi and the field-AP (nmcli) switching both require it.
User=root
WorkingDirectory=/home/pi/Big-Beautiful-Box
Environment=PYTHONUNBUFFERED=1
Environment=ROTORLINK_LOG=INFO
# Field hotspot: the Pi hosts its own WiFi AP (SSID = the trailer name) when no
# known network is in range, and joins a known network (e.g. the hangar) only when
# its signal clears the gate — see ROTORLINK_STA_JOIN_SIGNAL/STA_DROP_SIGNAL in
# rotorlink/config-less env (network_manager.py). ARMED fleet-wide. To DISABLE on a
# specific box: set the line below to 0 (or drop-in), then daemon-reload && restart.
Environment=ROTORLINK_AP_ENABLED=1
ExecStart=/usr/bin/python3 -m rotorlink
Restart=always
RestartSec=3
# Tunables (see rotorlink/config.py for the full list):
# Environment=ROTORLINK_WS_PORT=8765
# Environment=ROTORLINK_NAME=Trailer 1
# Environment=ROTORLINK_SERIAL=TS-01

[Install]
WantedBy=multi-user.target
