# deploy/bbb-logrotate.service
[Unit]
Description=Rotate BBB application logs

[Service]
Type=oneshot
ExecStart=/usr/sbin/logrotate -f /etc/logrotate.d/bbb
